# build/supervisord.conf
[supervisord]
nodaemon=true

[program:ntfy]
command=/usr/bin/ntfy serve
autostart=true
autorestart=true
stderr_logfile=/var/log/ntfy.err.log
stdout_logfile=/var/log/ntfy.out.log

[program:fastapi]
command=uvicorn --host=0.0.0.0 --port 8000 --app-dir /app  main:app
autostart=true
autorestart=true
stderr_logfile=/var/log/fastapi.err.log
stdout_logfile=/var/log/fastapi.out.log